# roles/millennium-falcon/templates/flexget.service.j2
[Unit]
Description=FlexGet Daemon
After=network.target

[Service]
ExecStart=/home/{{ millennium_falcon.pilot }}/.flexget/virtualenv/bin/flexget execute

[Install]
WantedBy=default.target
